# linux/howto.env.cron

0 5 * * * . $HOME/.profile; /path/to/command/to/run

#!/bin/bash
source $HOME/.bash_profile
some_other_cmd
